# deploy/rag-mcp-backfill.service
[Unit]
Description=rag-mcp incremental backfill (catches long-lived sessions before they end)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/projects/rag-mcp
ExecStart=%h/projects/rag-mcp/.venv/bin/rag-mcp backfill --source all
Environment=RAG_DISTILL_MODEL=Qwen3.6-35b-1M-P1-MTP-NGRAM
Nice=10
